# conf/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock                       ; path to your socket file

[supervisord]
logfile=/tmp/supervisord.log                    ; supervisord log file
logfile_maxbytes=50MB                           ; maximum size of logfile before rotation
logfile_backups=10                              ; number of backed up logfiles
loglevel=warn                                   ; info, debug, warn, trace
pidfile=/tmp/supervisord.pid                    ; pidfile location
nodaemon=false                                  ; run supervisord as a daemon
minfds=1024                                     ; number of startup file descriptors
minprocs=200                                    ; number of process descriptors
user=lepture                                    ; default user
childlogdir=/home/lepture/logs/

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock           ; use a unix:// URL  for a unix socket

[program:tornado-8011]
command=june/app.py --config=/home/lepture/config/june.config --port=8011 --logging=info --log_file_prefix=/home/lepture/logs/june-8011.log
stderr_logfile = /home/lepture/logs/tornado-8011-stderr.log
stdout_logfile = /home/lepture/logs/tornado-8011-stdout.log
